# agent/bin/deb/postinst
#!/bin/sh
set -e

SERVICE_NAME="iotistica-agent"
AGENT_DIR="/opt/iotistic/agent"
DATA_DIR="/var/lib/iotistic/agent"
LOG_DIR="/var/log/iotistic"
CONFIG_DIR="/etc/iotistic"
ENV_FILE="$CONFIG_DIR/agent.env"
SETUP_DONE="/var/lib/iotistic/setup-complete"
SETUP_SCRIPT="/usr/local/sbin/iotistica-setup.sh"
AGENT_SERVICE="/lib/systemd/system/${SERVICE_NAME}.service"

MQTT_PORT="${MQTT_BROKER_PORT:-1883}"
MQTT_USERNAME="${MQTT_USERNAME:-admin}"
MQTT_PASSWORD="${MQTT_PASSWORD:-admin}"

# ── User ──────────────────────────────────────────────────────────────────────

if ! id -u iotistic > /dev/null 2>&1; then
    useradd --system --home-dir /opt/iotistic --shell /bin/bash \
        --comment "Iotistica Agent" iotistic
fi

# ── Directories ───────────────────────────────────────────────────────────────

mkdir -p "$DATA_DIR" "$LOG_DIR" "$CONFIG_DIR"
mkdir -p /opt/iotistic/.config/node-opcua-default-nodejs/PKI
chown -R iotistic:iotistic /opt/iotistic/.config
chmod -R 755 /opt/iotistic/.config
chown -R iotistic:iotistic "$AGENT_DIR"
chown -R iotistic:iotistic "$DATA_DIR"
chown -R iotistic:iotistic "$LOG_DIR"

# ── Default config (fresh install only) ──────────────────────────────────────

if [ ! -f "$ENV_FILE" ]; then
    cat > "$ENV_FILE" << EOFENV
# Iotistica Agent configuration
# Edit this file and restart the service: systemctl restart iotistica-agent
#
# Leave IOTISTICA_API and PROVISIONING_KEY blank to run in standalone mode.
# Open the admin UI at http://<device-ip>:48484/admin/ to configure the device.

NODE_ENV=production
DEPLOYMENT_TYPE=systemd
DEVICE_API_PORT=48484
API_SECURITY_MODE=LOCAL_NETWORK
STANDALONE=true

# Cloud connectivity (optional — provision via admin UI instead)
IOTISTICA_API=
PROVISIONING_KEY=

# Logging
LOG_LEVEL=info
DATA_DIR=/var/lib/iotistic/agent
LOG_DIR=/var/log/iotistic

# Local Mosquitto MQTT broker (configured by installer)
MQTT_BROKER_URL=mqtt://localhost:${MQTT_PORT}
MQTT_USERNAME=${MQTT_USERNAME}
MQTT_PASSWORD=${MQTT_PASSWORD}
EOFENV
    chown iotistic:iotistic "$ENV_FILE"
    chmod 600 "$ENV_FILE"
fi

# ── Setup script ──────────────────────────────────────────────────────────────
# All apt-get/curl installs run here, outside the dpkg lock, via a
# systemd oneshot service started with --no-block after postinst exits.

cat > "$SETUP_SCRIPT" << 'EOFSETUP'
#!/bin/sh
set -e

SETUP_DONE="/var/lib/iotistic/setup-complete"
MQTT_PORT="@@MQTT_PORT@@"
MQTT_USERNAME="@@MQTT_USERNAME@@"
MQTT_PASSWORD="@@MQTT_PASSWORD@@"

if [ -f "$SETUP_DONE" ]; then
    echo "iotistica-setup: already complete — skipping"
    exit 0
fi

# dpkg holds its locks while postinst runs. We start with --no-block so this
# script runs in a separate process; wait for the locks to clear before
# calling apt-get (typically < 2 s after dpkg exits).
echo "iotistica-setup: waiting for package manager lock..."
i=0
while [ "$i" -lt 30 ]; do
    if flock --nonblock /var/lib/dpkg/lock-frontend true 2>/dev/null && \
       flock --nonblock /var/lib/dpkg/lock true 2>/dev/null; then
        break
    fi
    i=$((i + 1))
    sleep 2
done
echo "iotistica-setup: package manager available"

export DEBIAN_FRONTEND=noninteractive

# ── Node.js 22 ───────────────────────────────────────────────────────────────

NODE_OK=0
if command -v node > /dev/null 2>&1; then
    NODE_MAJOR=$(node -e "process.stdout.write(String(process.version.split('.')[0].slice(1)))" 2>/dev/null || echo "0")
    if [ "$NODE_MAJOR" -ge 22 ] 2>/dev/null; then
        NODE_OK=1
    fi
fi
if [ "$NODE_OK" != "1" ]; then
    echo "iotistica-setup: installing Node.js 22..."
    curl -fsSL https://deb.nodesource.com/setup_22.x | bash - > /dev/null 2>&1
    apt-get install -y -qq nodejs > /dev/null
    echo "iotistica-setup: ✓ Node.js $(node --version)"
fi

# ── Docker ────────────────────────────────────────────────────────────────────

if ! command -v docker > /dev/null 2>&1; then
    echo "iotistica-setup: installing Docker..."
    curl -fsSL https://get.docker.com -o /tmp/get-docker.sh
    sh /tmp/get-docker.sh > /dev/null 2>&1
    rm -f /tmp/get-docker.sh
    systemctl enable docker > /dev/null 2>&1
    systemctl start docker > /dev/null 2>&1 || true
    echo "iotistica-setup: ✓ Docker installed"
fi
if getent group docker > /dev/null 2>&1; then
    usermod -aG docker iotistic
fi

# ── Mosquitto ─────────────────────────────────────────────────────────────────

if ! command -v mosquitto > /dev/null 2>&1; then
    echo "iotistica-setup: installing Mosquitto..."
    apt-get install -y -qq --no-install-recommends mosquitto mosquitto-clients > /dev/null
    echo "iotistica-setup: ✓ Mosquitto installed"
fi
if ! command -v setfacl > /dev/null 2>&1; then
    apt-get install -y -qq --no-install-recommends acl > /dev/null 2>&1 || true
fi

echo "iotistica-setup: configuring Mosquitto..."
mkdir -p /etc/mosquitto/conf.d /var/log/mosquitto
chgrp mosquitto /etc/mosquitto 2>/dev/null || true
chmod 2755 /etc/mosquitto 2>/dev/null || true

if ! grep -q 'include_dir /etc/mosquitto/conf.d' /etc/mosquitto/mosquitto.conf 2>/dev/null; then
    echo "include_dir /etc/mosquitto/conf.d" >> /etc/mosquitto/mosquitto.conf
fi

# Auto-detect port conflict and switch to an alternative if needed
if ss -tlnp 2>/dev/null | grep -q ":${MQTT_PORT}[[:space:]]"; then
    echo "iotistica-setup: port $MQTT_PORT is already in use"
    SWITCHED=0
    for TRY_PORT in 8883 1884 1885 1886; do
        if ! ss -tlnp 2>/dev/null | grep -q ":${TRY_PORT}[[:space:]]"; then
            MQTT_PORT="$TRY_PORT"
            SWITCHED=1
            echo "iotistica-setup: auto-switching Mosquitto to port $MQTT_PORT"
            sed -i "s|MQTT_BROKER_URL=mqtt://localhost:[0-9]*|MQTT_BROKER_URL=mqtt://localhost:${MQTT_PORT}|" /etc/iotistic/agent.env
            break
        fi
    done
    if [ "$SWITCHED" = "0" ]; then
        echo "iotistica-setup: WARNING all candidate MQTT ports are in use — Mosquitto may not start"
    fi
fi

cat > /etc/mosquitto/conf.d/iotistica.conf << EOFMOSQ
# Iotistica managed configuration

listener $MQTT_PORT
allow_anonymous false

password_file /etc/mosquitto/passwd
acl_file /etc/mosquitto/acl
EOFMOSQ

touch /etc/mosquitto/passwd /etc/mosquitto/acl
touch /var/log/mosquitto/mosquitto.log

mosquitto_passwd -b -c /etc/mosquitto/passwd "$MQTT_USERNAME" "$MQTT_PASSWORD"

cat > /etc/mosquitto/acl << EOFACL
# Iotistica managed ACL
user $MQTT_USERNAME
topic readwrite #
EOFACL

chown root:mosquitto /etc/mosquitto/passwd /etc/mosquitto/acl
chmod 0640 /etc/mosquitto/passwd /etc/mosquitto/acl
chown mosquitto:mosquitto /var/log/mosquitto/mosquitto.log
chmod 640 /var/log/mosquitto/mosquitto.log
usermod -aG mosquitto iotistic 2>/dev/null || true

if command -v setfacl > /dev/null 2>&1; then
    setfacl -m u:iotistic:rwx /etc/mosquitto || true
    setfacl -m u:iotistic:rw /etc/mosquitto/passwd || true
    setfacl -m u:iotistic:rw /etc/mosquitto/acl || true
fi

cat > /usr/local/bin/iotistica-mqtt-reload.sh << 'EOFRELOAD'
#!/bin/sh
set -e
chown root:mosquitto /etc/mosquitto/passwd /etc/mosquitto/acl
chmod 0640 /etc/mosquitto/passwd /etc/mosquitto/acl
if command -v setfacl > /dev/null 2>&1; then
    setfacl -m u:iotistic:rw /etc/mosquitto/passwd || true
    setfacl -m u:iotistic:rw /etc/mosquitto/acl || true
fi
systemctl reload mosquitto
EOFRELOAD
chmod +x /usr/local/bin/iotistica-mqtt-reload.sh

cat > /etc/sudoers.d/iotistica-mqtt-reload << 'EOFSUDO'
iotistic ALL=(root) NOPASSWD: /usr/local/bin/iotistica-mqtt-reload.sh
EOFSUDO
chmod 440 /etc/sudoers.d/iotistica-mqtt-reload

cat > /etc/systemd/system/iotistica-mqtt-reload.service << 'EOFMQSVC'
[Unit]
Description=Reload Mosquitto MQTT auth files after credential update
After=mosquitto.service
Requires=mosquitto.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/iotistica-mqtt-reload.sh
EOFMQSVC

cat > /etc/systemd/system/iotistica-mqtt-reload.path << 'EOFMQPATH'
[Unit]
Description=Watch Mosquitto passwd file for credential updates

[Path]
PathModified=/etc/mosquitto/passwd
Unit=iotistica-mqtt-reload.service

[Install]
WantedBy=multi-user.target
EOFMQPATH

systemctl daemon-reload
systemctl enable mosquitto > /dev/null 2>&1
systemctl enable --now iotistica-mqtt-reload.path > /dev/null 2>&1
if systemctl restart mosquitto; then
    echo "iotistica-setup: ✓ Mosquitto running on port $MQTT_PORT"
else
    echo "iotistica-setup: WARNING Mosquitto failed to start on port $MQTT_PORT"
    echo "  Diagnosis:  journalctl -u mosquitto -n 20 --no-pager"
    echo "  Port check: ss -tlnp | grep :$MQTT_PORT"
    echo "  After fix:  sudo systemctl restart mosquitto && sudo systemctl restart iotistica-agent"
fi

# ── Mark complete and start agent ─────────────────────────────────────────────

touch "$SETUP_DONE"
echo "iotistica-setup: ✓ setup complete — starting Iotistica Agent"
systemctl enable iotistica-agent > /dev/null 2>&1
systemctl start iotistica-agent || true

echo ""
echo "✓ Iotistica Agent is running."
echo "  Admin UI:  http://localhost:48484/admin/"
echo "  MQTT:      localhost:$MQTT_PORT  (user: $MQTT_USERNAME)"
echo "  Status:    systemctl status iotistica-agent"
echo "  Logs:      journalctl -u iotistica-agent -f"
echo ""
EOFSETUP

chmod +x "$SETUP_SCRIPT"

# Inject MQTT vars (placeholders used because the outer heredoc is quoted)
sed -i \
    -e "s|@@MQTT_PORT@@|${MQTT_PORT}|g" \
    -e "s|@@MQTT_USERNAME@@|${MQTT_USERNAME}|g" \
    -e "s|@@MQTT_PASSWORD@@|${MQTT_PASSWORD}|g" \
    "$SETUP_SCRIPT"

# ── systemd setup service ─────────────────────────────────────────────────────

cat > /etc/systemd/system/iotistica-setup.service << 'EOFSETUPSVC'
[Unit]
Description=Iotistica Agent First-Run Setup
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/iotistic/setup-complete

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/iotistica-setup.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=iotistica-setup
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
EOFSETUPSVC

# ── systemd agent service ─────────────────────────────────────────────────────

if ! command -v systemctl > /dev/null 2>&1; then
    echo "Warning: systemd not found — skipping service setup."
    echo "Start the agent manually: node $AGENT_DIR/dist/app.js"
    exit 0
fi

# NodeSource installs node to /usr/bin/node; use that as the canonical path.
# If node is already present (upgrade), use the actual path.
NODE_BIN=$(command -v node 2>/dev/null || echo "/usr/bin/node")

if [ -f "$AGENT_DIR/dist/app.js" ]; then
    APP_JS="$AGENT_DIR/dist/app.js"
elif [ -f "$AGENT_DIR/app.js" ]; then
    APP_JS="$AGENT_DIR/app.js"
elif [ -f "$AGENT_DIR/dist/src/app.js" ]; then
    APP_JS="$AGENT_DIR/dist/src/app.js"
else
    APP_JS="$AGENT_DIR/dist/app.js"
fi

cat > "$AGENT_SERVICE" << EOF
[Unit]
Description=Iotistica Agent - IoT Device Management Service
Documentation=https://docs.iotistica.com
After=network-online.target iotistica-setup.service docker.service mosquitto.service
Wants=network-online.target iotistica-setup.service
Requires=docker.service
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=notify
User=iotistic
Group=iotistic
WorkingDirectory=$AGENT_DIR
EnvironmentFile=$ENV_FILE
Environment=NODE_ENV=production
Environment=DEPLOYMENT_TYPE=systemd
ExecStart=$NODE_BIN $APP_JS
Restart=always
RestartSec=5
WatchdogSec=60
NotifyAccess=all
TimeoutStartSec=120
TimeoutStopSec=20
KillMode=control-group
StandardOutput=journal
StandardError=journal
SyslogIdentifier=$SERVICE_NAME
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/iotistic /var/log/iotistic $AGENT_DIR /var/run/docker.sock
LockPersonality=true
MemoryAccounting=true
CPUAccounting=true
LimitNOFILE=65536
LimitNPROC=65536
MemoryHigh=250M
MemoryMax=300M
TasksMax=512
CPUQuota=80%

[Install]
WantedBy=multi-user.target
EOF

# ── Enable and kick off ───────────────────────────────────────────────────────

systemctl daemon-reload
systemctl enable iotistica-setup.service > /dev/null 2>&1
systemctl enable "$SERVICE_NAME" > /dev/null 2>&1

if [ -f "$SETUP_DONE" ]; then
    # Upgrade: dependencies already installed — just restart the agent
    systemctl restart "$SERVICE_NAME" || true
    echo ""
    echo "✓ Iotistica Agent updated and restarted."
    echo "  Status: systemctl status $SERVICE_NAME"
    echo ""
else
    # Fresh install: start setup in background so dpkg can release its locks first
    systemctl start --no-block iotistica-setup.service 2>/dev/null || true
    echo ""
    echo "✓ Iotistica Agent package installed."
    echo "  Installing Node.js, Docker, and Mosquitto in the background..."
    echo "  Setup progress: journalctl -u iotistica-setup -f"
    echo "  The agent starts automatically when setup completes (~1-3 min)."
    echo ""
fi
